# zeusci/conf/citemplate/config/supervisord.conf
[unix_http_server]
file=%(here)s/../var/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=%(here)s/../var/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://%(here)s/../var/supervisor.sock ; use a unix:// URL  for a unix socket


; ===========
;  webserver
; ===========

[program:webserver]
command=%(here)s/../venv/bin/python %(here)s/../app/manage.py runserver 127.0.0.1:23115
stdout_logfile=%(here)s/../log/webserver.log
stderr_logfile=%(here)s/../log/webserver.error.log
stopsignal=KILL
killasgroup=true
stopasgroup=true


; ========
;  worker
; ========

[program:worker]
command=%(here)s/../venv/bin/python %(here)s/../app/manage.py celery worker --loglevel=INFO

numprocs=1
stdout_logfile=%(here)s/../log/celery.log
stderr_logfile=%(here)s/../log/celery.error.log
autostart=true
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

; When resorting to send SIGKILL to the program to terminate it
; send SIGKILL to its whole process group instead,
; taking care of its children as well.
killasgroup=true
